# misc/gpu-homeassistant.timer
[Unit]
Description=Send GPU temperature to homeassistant

[Timer]
OnBootSec=5min
OnUnitActiveSec=15sec

[Install]
WantedBy=timers.target

